# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=1

[include]
files=/tmp/.docker_generated/supervisord/*.conf
